# build/debian/lib/systemd/system/fastlane-sidecar.service
[Unit]
Description=Fastlane Sidecar - MEV capturing sidecar for Monad Validator
After=network.target
Wants=network.target

[Service]
Type=simple
User=fastlane
Group=fastlane
ExecStart=/usr/bin/fastlane-sidecar -log-level=info -http-port=8080
Restart=always

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true

RuntimeDirectory=fastlane-sidecar

[Install]
WantedBy=multi-user.target
