# debian/sidecar/fastlane-sidecar.service
[Unit]
Description=Fastlane Sidecar - MEV capturing sidecar for Monad Validator
After=network.target
Wants=network.target

[Service]
Type=simple
User=monad
Group=monad
EnvironmentFile=/home/monad/fastlane/.env
ExecStart=/usr/bin/fastlane-sidecar \
    --log-level info \
    --network "${NETWORK}"
Restart=always

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict

RuntimeDirectory=fastlane-sidecar

[Install]
WantedBy=multi-user.target
