# debian/sidecar/postinst
#!/bin/bash

set -e

# Reload systemd to recognize the new service
systemctl daemon-reload

echo "Fastlane Sidecar installed successfully!"
echo ""
echo "Note: The service runs as the 'monad' user to access the Monad validator's IPC socket."
echo ""
echo "Next steps:"
echo "  1. Configure if needed: sudo systemctl edit fastlane-sidecar"
echo "  2. Enable on boot: sudo systemctl enable fastlane-sidecar"
echo "  3. Start the service: sudo systemctl start fastlane-sidecar"
