# debian/postrm
#!/bin/bash

set -e

case "$1" in
    remove|purge)
        # Stop service if it's running
        if systemctl is-active --quiet fastlane-sidecar.service; then
            systemctl stop fastlane-sidecar.service || true
        fi

        # Disable service if it's enabled
        if systemctl is-enabled --quiet fastlane-sidecar.service; then
            systemctl disable fastlane-sidecar.service || true
        fi

        # Reload daemon after service operations
        systemctl daemon-reload || true

        # Remove user and group
        if getent passwd fastlane > /dev/null 2>&1; then
            userdel fastlane
        fi

        if getent group fastlane > /dev/null 2>&1; then
            groupdel fastlane
        fi

        ;;
esac
