# debian/postrm
#!/bin/bash

set -e

case "$1" in
    remove|purge)
        # Stop and disable service
        systemctl stop fastlane-sidecar.service || true
        systemctl disable fastlane-sidecar.service || true
        systemctl daemon-reload

        # Remove user and group
        if getent passwd fastlane > /dev/null 2>&1; then
            userdel fastlane
        fi

        if getent group fastlane > /dev/null 2>&1; then
            groupdel fastlane
        fi

        ;;
esac
